# ndi-monitor.service
[Unit]
Description=NDI Monitor Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nlc-atlantic-media
WorkingDirectory=/home/nlc-atlantic-media/ndi-monitor
Environment="PATH=/home/nlc-atlantic-media/ndi-monitor/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NDI_LIB_PATH=/usr/local/lib/libndi.so"
ExecStart=/home/nlc-atlantic-media/ndi-monitor/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
